# install/supervisord.ini
[supervisord]
nodaemon=true
logfile=/logs/supervisord-HHR.log
logfile_maxbytes=50MB
logfile_backups=4
loglevel=info
pidfile=/var/run/supervisord.pid

[program:conf-watch]
command=node /opt/conf-watch/confd-watch -h etcd
autorestart=true
stdout_logfile=/logs/conf-watch.log
stdout_logfile_maxbytes=2MB
stdout_logfile_backups=7

[program:conf-register]
command=node /opt/conf-register/conf-register -h etcd -c /config/config.json
autorestart=true
stdout_logfile=/logs/conf-register.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=7

[program:physiodom-app]
directory=/src/
command=node server/server.js -c /config/config.json
autostart=true
autorestart=true
stdout_logfile=/logs/app.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=7
stderr_logfile=/logs/app-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=7

[program:nginx]
command=nginx
autorestart=true